# deploy/hourly-update.cron.example
# Run the production update once per hour.
# Replace /path/to/music-review with the repository path on the server.
# New reviews trigger metadata enrichment and artist image prefetch for that batch.
0 * * * * cd /path/to/music-review && mkdir -p logs && docker compose run --rm music-review-update >> logs/hourly-update.log 2>&1

# Prefetch artist images overnight (cache + local JPGs).
0 2 * * * cd /path/to/music-review && mkdir -p logs && ARTIST_IMAGE_DOWNLOAD=true ARTIST_IMAGE_RESOLVE_ON_DEMAND=false hatch run artist-image-batch --missing-only --all -v >> logs/artist-image-batch.log 2>&1
